# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/home/user/run/supervisord.pid
loglevel=info

; ---------- GitHub Sync Daemon ----------
[program:syncd]
command=/bin/bash -lc 'umask 000 && python3 /home/user/sync_to_github.py'
directory=/home/user
autostart=true
autorestart=true
startsecs=2
priority=10
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---------- Stardew Valley Server ----------
[program:stardew]
command=/bin/bash -lc 'umask 000 && /home/user/scripts/wait-sync-ready.sh; chmod -R 777 /home/steam/.sdv-backup || true; /home/steam/entrypoint.sh'
directory=/home/steam
autostart=true
autorestart=true
startsecs=5
priority=20
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=180
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---------- FRP Client (switchable) ----------
[program:frpc]
command=/bin/bash -c "/home/user/frp/frp-entry.sh"
directory=/home/user/frp
autostart=true
autorestart=true
startsecs=3
priority=25
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---------- Web Monitor ----------
[program:webmonitor]
command=python3 /home/user/web/app.py
directory=/home/user/web
autostart=true
autorestart=true
startsecs=3
priority=30
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
