# supervisor/supervisord.conf
; Supervisord for sdv-hf (aligned with mc-hf layout)
[supervisord]
logfile=/dev/null
pidfile=/home/user/supervisord.pid
nodaemon=true
strip_ansi=false
loglevel=info
childlogdir=/home/user/logs

; ---------- GitHub Sync Daemon ----------
; Start sync web/API at :5321 if available. If GITHUB_REPO/PAT not set or
; sync module unavailable, keep running idle so other services are gated.
[program:syncd]
command=/bin/bash -lc '
  if [[ -z "${GITHUB_REPO:-}" || -z "${GITHUB_PAT:-}" ]]; then
    echo "[syncd] GITHUB_REPO/PAT not set. Staying idle; gating others." >&2
    exec bash -lc "sleep infinity"
  fi
  if python3 -c "import sync" 2>/dev/null; then
    exec python3 -m sync
  else
    echo "[syncd] Python module 'sync' not found. Install/copy it, or mount it."
    exec bash -lc "sleep infinity"
  fi'
directory=/home/user
autostart=true
autorestart=true
startsecs=1
priority=15
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; ---------- Nginx (OpenResty) ----------
[program:nginx]
command=/bin/bash -lc '
  if command -v /usr/local/openresty/bin/openresty >/dev/null 2>&1; then
    exec /usr/local/openresty/bin/openresty -p /home/user/nginx -c /home/user/nginx/nginx.conf -g "daemon off;"
  elif command -v openresty >/dev/null 2>&1; then
    exec openresty -p /home/user/nginx -c /home/user/nginx/nginx.conf -g "daemon off;"
  elif command -v nginx >/dev/null 2>&1; then
    exec nginx -c /home/user/nginx/nginx.conf -g "daemon off;"
  else
    echo "[nginx] OpenResty/Nginx not installed; sleeping." >&2
    exec bash -lc "sleep infinity"
  fi'
directory=/home/user/nginx
autostart=true
autorestart=true
priority=40
stopsignal=TERM
stdout_logfile=/home/user/logs/nginx.supervisor.out
stderr_logfile=/home/user/logs/nginx.supervisor.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; ---------- Stardew Valley Server (gated by GitHub sync) ----------
[program:stardew]
command=/bin/bash -lc '
  if [[ -z "${GITHUB_REPO:-}" || -z "${GITHUB_PAT:-}" ]]; then
    echo "[stardew] Sync not configured; not starting." >&2
    exec bash -lc "sleep infinity"
  fi
  umask 000
  /home/user/scripts/wait-sync-ready.sh
  chmod -R 777 /home/steam/.sdv-backup || true
  exec /home/steam/entrypoint.sh'
directory=/home/steam
autostart=true
autorestart=true
startsecs=5
priority=20
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=180
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---------- FRP Client (gated by GitHub sync) ----------
[program:frpc]
command=/bin/bash -lc '
  if [[ -z "${GITHUB_REPO:-}" || -z "${GITHUB_PAT:-}" ]]; then
    echo "[frpc] Sync not configured; not starting." >&2
    exec bash -lc "sleep infinity"
  fi
  HIST_DIR="${HIST_DIR:-/home/steam/.sdv-backup}"; BRANCH="${GIT_BRANCH:-main}"
  if command -v git >/dev/null 2>&1; then
    echo "[frpc] Waiting for git HEAD to match origin/${BRANCH} in ${HIST_DIR}..."
    i=0
    until git -C "$HIST_DIR" rev-parse --git-dir >/dev/null 2>&1 \
      && git -C "$HIST_DIR" rev-parse "origin/${BRANCH}" >/dev/null 2>&1 \
      && [ "$(git -C "$HIST_DIR" rev-parse HEAD 2>/dev/null)" = "$(git -C "$HIST_DIR" rev-parse "origin/${BRANCH}" 2>/dev/null)" ]; do
      if (( i % 5 == 0 )); then
        L=$(git -C "$HIST_DIR" rev-parse HEAD 2>/dev/null || echo NA)
        O=$(git -C "$HIST_DIR" rev-parse "origin/${BRANCH}" 2>/dev/null || echo NA)
        echo "[frpc] waiting... local=${L} origin=${O}"
      fi
      i=$((i+1))
      sleep 2
    done
  else
    echo "[frpc] git not available; cannot verify sync. Not starting." >&2
    exec bash -lc "sleep infinity"
  fi
  exec /home/user/frp/frp-entry.sh'
directory=/home/user/frp
autostart=true
autorestart=true
startsecs=3
priority=25
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---------- File Browser (gated by GitHub sync) ----------
[program:filebrowser]
command=/bin/bash -lc '
  if [[ -z "${GITHUB_REPO:-}" || -z "${GITHUB_PAT:-}" ]]; then
    echo "[filebrowser] Sync not configured; not starting." >&2
    exec bash -lc "sleep infinity"
  fi
  HIST_DIR="${HIST_DIR:-/home/steam/.sdv-backup}"; BRANCH="${GIT_BRANCH:-main}"
  if command -v git >/dev/null 2>&1; then
    echo "[filebrowser] Waiting for git HEAD to match origin/${BRANCH} in ${HIST_DIR}..."
    i=0
    until git -C "$HIST_DIR" rev-parse --git-dir >/dev/null 2>&1 \
      && git -C "$HIST_DIR" rev-parse "origin/${BRANCH}" >/dev/null 2>&1 \
      && [ "$(git -C "$HIST_DIR" rev-parse HEAD 2>/dev/null)" = "$(git -C "$HIST_DIR" rev-parse "origin/${BRANCH}" 2>/dev/null)" ]; do
      if (( i % 5 == 0 )); then
        L=$(git -C "$HIST_DIR" rev-parse HEAD 2>/dev/null || echo NA)
        O=$(git -C "$HIST_DIR" rev-parse "origin/${BRANCH}" 2>/dev/null || echo NA)
        echo "[filebrowser] waiting... local=${L} origin=${O}"
      fi
      i=$((i+1))
      sleep 2
    done
  else
    echo "[filebrowser] git not available; cannot verify sync. Not starting." >&2
    exec bash -lc "sleep infinity"
  fi
  exec /home/user/filebrowser --address 0.0.0.0 --port 8000 --root / --database /home/user/filebrowser.db'
directory=/home/user
autostart=true
autorestart=true
startsecs=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---------- GoTTY Web Terminal ----------
[program:gotty]
command=/bin/bash -lc '/home/user/gotty --address 127.0.0.1 --port 8080 --path /terminal -w bash'
directory=/home/user
autostart=true
autorestart=true
startsecs=2
priority=28
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

